# homeassistant_gitops/rootfs/etc/services.d/ha_gitops_bridge/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

PORT="${INGRESS_PORT:-8099}"
exec uv run --project /app python /app/main.py --port "${PORT}"
