# homeassistant_gitops/rootfs/etc/services.d/homeassistant_gitops/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the Home Assistant GitOps Bridge service
# ==============================================================================
set -e

exec python /app/main.py --port 8099
